# meta-lmp-base/recipes-support/fioconfig/fioconfig/fioconfig.path
[Unit]
Description=Foundries.io configuration management path monitor

[Path]
PathExists=/var/sota/sota.toml

[Install]
WantedBy=multi-user.target
