# meta-lmp-base/recipes-support/ostree-pending-reboot/ostree-pending-reboot/ostree-pending-reboot.timer.in
[Unit]
Description=Automatic OSTree Update Reboot Scheduling

[Timer]
OnCalendar=*:0/@@OSTREE_PENDING_REBOOT_CHECK_MINUTES@@

[Install]
WantedBy=multi-user.target
