# meta-lmp-base/recipes-support/docker-auto-prune/docker-auto-prune/docker-auto-prune.service
[Unit]
Description=Automatic Docker Prune
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune --force
