# meta-lmp-base/recipes-sota/aktualizr/aktualizr/aktualizr-lite.path
[Unit]
Description=Aktualizr Lite SOTA Client path monitor
ConditionPathExists=/etc/systemd/system/multi-user.target.wants/aktualizr-lite.service
ConditionPathExists=!/usr/bin/mbedCloudClient

[Path]
PathExists=/var/sota/sota.toml

[Install]
WantedBy=multi-user.target
