# meta-lmp-base/recipes-sota/aktualizr/aktualizr/aktualizr-lite.service.in
[Unit]
Description=Aktualizr Lite SOTA Client
After=network.target boot-complete.target
Requires=boot-complete.target
ConditionPathExists=|/var/sota/sota.toml
# this should be StartLimitBurst * RestartSec + 1
StartLimitIntervalSec=541
StartLimitBurst=3
StartLimitAction=reboot

[Service]
User=root
RestartSec=180
Restart=always
ExecStartPre=/usr/bin/mkdir -p /run/aktualizr
Environment="TMPDIR=/run/aktualizr"
Environment="COMPOSE_HTTP_TIMEOUT=@@COMPOSE_HTTP_TIMEOUT@@"
Environment="REGISTRY_AUTH_FILE=@@DOCKER_CRED_HELPER_CFG@@"
# Allowed values are in the [1, 10] range. Aklite adjusts value to the corresponding limit if the specified value is higher or lower than [1, 10].
ExecStart=/usr/bin/aktualizr-lite daemon

[Install]
WantedBy=multi-user.target
