# meta-lmp-base/recipes-sota/aktualizr/aktualizr/aktualizr-lite.service.in
[Unit]
Description=Aktualizr Lite SOTA Client
After=network.target boot-complete.target
Requires=boot-complete.target
ConditionPathExists=|/var/sota/sota.toml
ConditionPathExists=|/usr/lib/sota/conf.d/10-lite-public-stream.toml
ConditionPathExists=!/usr/bin/mbedCloudClient

[Service]
RestartSec=180
Restart=always
ExecStartPre=/usr/bin/mkdir -p /run/aktualizr
Environment="TMPDIR=/run/aktualizr"
Environment="COMPOSE_HTTP_TIMEOUT=@@COMPOSE_HTTP_TIMEOUT@@"
ExecStart=/usr/bin/aktualizr-lite daemon

[Install]
WantedBy=multi-user.target
