# meta-lmp-base/recipes-sota/aktualizr/aktualizr/aktualizr-lite.path
[Unit]
Description=Aktualizr Lite SOTA Client path monitor
After=network.target
ConditionPathExists=!/usr/bin/mbedCloudClient

[Path]
PathExists=/var/sota/sota.toml
PathExists=/usr/lib/sota/conf.d/10-lite-public-stream.toml

[Install]
WantedBy=multi-user.target
