# meta-lmp-base/recipes-sota/aktualizr/aktualizr/aktualizr-lite.path
[Unit]
Description=Aktualizr Lite SOTA Client path monitor
ConditionPathExists=!/usr/bin/mbedCloudClient

[Path]
PathExists=/var/sota/sota.toml

[Install]
WantedBy=multi-user.target
