# recipes-support/docker-auto-prune/docker-auto-prune/docker-auto-prune.timer.in
[Unit]
Description=Automatic Docker Prune Timer

[Timer]
OnCalendar=@@DOCKER_PRUNE_ONCALENDAR@@
Persistent=false

[Install]
WantedBy=timers.target
