# meta-lmp-base/recipes-support/compose-apps-early-start/compose-apps-early-start/compose-apps-early-start-recovery.service
[Unit]
Description=Recovery for compose-apps-early-start.service

[Service]
Type=oneshot
ExecStart=/usr/bin/compose-apps-early-start-recovery
