# recipes-support/lmp-auto-hostname/lmp-auto-hostname/lmp-auto-hostname.service.in
[Unit]
Description=Linux microPlatform Auto Hostname Update
DefaultDependencies=no
Before=network-pre.target avahi-daemon.service
After=local-fs.target dbus.service

[Service]
Environment=MACHINE=@@LMP_HOSTNAME_MACHINE@@ MODE=@@LMP_HOSTNAME_MODE@@ NETDEVICE=@@LMP_HOSTNAME_NETDEVICE@@
ExecStart=/usr/bin/lmp-update-hostname
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
